# images/capi/ansible/roles/kubernetes/templates/usr/lib/systemd/system/kubelet.service.d/11-resource-sizing.conf
[Service]
EnvironmentFile=-/run/kubelet/extra-args.env
ExecStartPre={{ sysusr_prefix }}/libexec/kubernetes/kubelet-resource-sizing.sh
