# images/capi/ansible/roles/providers/files/etc/systemd/system/cloud-final.service.d/boot-order.conf
[Unit]
After=containerd.service
Wants=containerd.service