# images/capi/ansible/roles/containerd/templates/etc/systemd/system/containerd.service.d/limit-nofile.conf
[Service]
# LimitNOFILE=infinity on Mariner means 1073741816, which has caused issues
# running some software in containers such as mysql5 and sshd.
LimitNOFILE=1048576
