# images/capi/ansible/roles/containerd/templates/etc/systemd/system/containerd.service.d/max-tasks.conf
[Service]
# Do not limit the number of tasks that can be spawned by containerd
TasksMax=infinity
